# .config/systemd/user/mailsync.service
[Unit]
Description=mailsync
Wants=mailsync.timer

[Service]
Type=oneshot
ExecStart=/home/kgn/bin/mailsync

[Install]
WantedBy=default.target

