# .config/systemd/user/sink.service
[Unit]
Description=sink
Wants=sink.timer

[Service]
Type=oneshot
ExecStart=/home/kgn/bin/sinks.sh

[Install]
WantedBy=default.target

