# .config/systemd/user/taskwarrior.service
[Unit]
Description=Taskwarrior
Wants=taskwarrior.timer

[Service]
Type=oneshot
ExecStart=/sbin/task sync

[Install]
WantedBy=default.target

